# supervisord.conf
[program:schools_app]
environment: %(host_envs)s
command: /usr/bin/python -m schools_app
directory: /home/ubuntu
autostart: true