# supervisord.conf
[program:schools_app]
command: /usr/bin/python -m schools_app
directory: /home/ubuntu
autostart: true