# supervisord.conf
[program:schools_app]
command: /usr/bin/python -m bookapp
directory: /home/ubuntu
autostart: true